# Cubietruck/usr/lib/systemd/system/cubietruck_leds_off.service
[Unit]
Description=Disable the LEDs of Cubietruck
[Service]
ExecStart=/usr/bin/cubietruck_leds_off.sh
[Install]
WantedBy=multi-user.target
